# _sh/crontab
####
###
##
# ADD this line to the END of /etc/crontab
@reboot root bash /srv/fe/_startup/start.sh